# services/php-cli/crontabs/root
*/1 * * * * echo 'crontab'  > /dev/null
